# config/hadoop_config/node02/etc/profile.d/hbase.sh
if [ -z $HBASE_HOME ]; then
    export HBASE_HOME=/app/hbase-1.3.1
    export PATH=$PATH:$HBASE_HOME/bin
fi
