# config/hadoop_config/node01/etc/profile.d/hadoop.sh
if [ -z "$HADOOP_HOME" ]; then
    export HADOOP_HOME=/app/hadoop-2.7.1
    export PATH=$PATH:$HADOOP_HOME/bin
fi
