# config/hadoop_config/node01/etc/profile.d/hive.sh
if [ -z $HIVE_HOME ]; then
    export HIVE_HOME=/app/apache-hive-2.1.1-bin
    export PATH=$PATH:$HIVE_HOME/bin
fi
